# ChristmasSocks/systemd/socks.service
[Unit]
Description=ChristmasSocks Remote Server Management Software

[Service]
Type=simple
Restart=always
RestartSec=1
User=username
ExecStart=/usr/bin/socks

[Install]
WantedBy=multi-user.target
